# misc/pkg/scripts/indns.postinstall
#!/bin/sh
# RPM %post / deb postinst for innerstack-indns
systemctl daemon-reload >/dev/null 2>&1 || true
if [ "$1" = "1" ] || [ "$1" = "configure" ]; then
    systemctl preset indns.service >/dev/null 2>&1 || true
fi
