# misc/pkg/scripts/indns.preremove
#!/bin/sh
# RPM %preun / deb prerm for innerstack-indns
if [ "$1" = "0" ] || [ "$1" = "remove" ]; then
    systemctl --no-warn stop indns.service >/dev/null 2>&1 || true
    systemctl --no-warn disable indns.service >/dev/null 2>&1 || true
fi
