# misc/pkg/scripts/ingate.postinstall
#!/bin/sh
# RPM %post / deb postinst for innerstack-ingate
systemctl daemon-reload >/dev/null 2>&1 || true
if [ "$1" = "1" ] || [ "$1" = "configure" ]; then
    systemctl preset ingate.service >/dev/null 2>&1 || true
fi
